# code_snippets/timelapse/etc/systemd/system/timelapse.service
[Unit]
Description=Nx snaphot timelapse
[Service]
ExecStart=/opt/nxos/timelapse.sh